# olds/.local/etc/anacrontab
#period in days   delay in minutes   job-identifier   command
1 5 1 /usr/bin/trash-put /home/datsfilipe/.local/share/nvim/backup/*
